# .config/systemd/user/20-20-20.service
[Unit]
Description=20-20-20 rule notification

[Service]
Type=oneshot
ExecStart=notify-send "look into the distance for 20 secs"
