# .config/systemd/user/go-to-sleep.service
[Unit]
Description=go-to-sleep notification

[Service]
Type=oneshot
ExecStartPre=notify-send "GO TO SLEEP"
ExecStart=shutdown
