# .config/systemd/user/go-to-sleep.timer
[Unit]
Description=go-to-sleep timer

[Timer]
OnCalendar=*-*-* 00:00:00

[Install]
WantedBy=timers.target
